# example/rootfs/etc/crontabs/echo-test
*/1 * * * * /etc/service.d/cron/echo-test
